# configs/volcanogunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/volcanogunicorn.sock

[Install]
WantedBy=sockets.target
